# system_files/usr/lib/systemd/user/spinosauros-first-login.service
[Unit]
Description=Spinosauros first login setup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/first-login.sh

[Install]
WantedBy=default.target
