# system_files/etc/xdg/autostart/appgrid.desktop
[Desktop Entry]
Type=Application
Name=Configure AppGrid
Exec=/usr/share/ublue-os-spinosauros/appgrid.sh
OnlyShowIn=GNOME;
NoDisplay=true
X-GNOME-Autostart-enabled=true
